# systemd.d/lunar-rebuild.service
[Unit]
Description=Lunar ISO rebuild service, for use after fixing thing manually

[Service]
ExecStart=/home/dagbrown/bin/lunar-rebuild
