# systemd.d/lunar-build.service
[Unit]
Description=Lunar ISO build service

[Service]
ExecStart=/home/dagbrown/bin/lunar-build
